# operations/systemd/ipo-ingest.service.d/environment.conf
# Reference copy of /etc/systemd/system/ipo-ingest.service.d/environment.conf.
# ALREADY APPLIED on the box (operations/README.md previously listed this as an outstanding step).
# Supplies UPSTOX_TOKEN for the new-IPO onset context pull; systemd reads the file as root.
[Service]
EnvironmentFile=/etc/ipo/context.env
