# operations/systemd/ipo-ingest.service.d/onfailure.conf
# Reference copy of /etc/systemd/system/ipo-ingest.service.d/onfailure.conf.
#
# Fires a Telegram warning when the ingest unit genuinely FAILS. The template message is fixed --
# "WARN: <unit> errored (systemd OnFailure)" -- and its only variable is the unit name.
#
# v3-DP DP-1 DELIBERATELY DOES NOT USE THIS for the recorder. See ./README.md: the message cannot
# distinguish "recorder leg failed, ingest fine" from "ingest failed", so overloading it would make
# this alert contradict the digest. The recorder reports through its own Recorder health row.
[Unit]
OnFailure=ipo-telegram-alert@%n.service
